# rootfs/etc/services.d/frigateconnect/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "Frigate Connect stopped."
